# deploy/systemd/gridpulse-dashboard.service
[Unit]
Description=GridPulse Frontend (Next.js)
After=network.target

[Service]
Type=simple
User=gridpulse
WorkingDirectory=/opt/gridpulse
EnvironmentFile=/etc/gridpulse/gridpulse.env
Environment=PORT=3000
Environment=HOSTNAME=0.0.0.0
ExecStart=/usr/bin/npm --prefix /opt/gridpulse/frontend run start
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
